# cyber-cloud-project/code/startup/run-startup.service
[Unit]
After=network.service

[Service]
ExecStart=/root/startup.sh

[Install]
WantedBy=default.target